# Library/LaunchAgents/gnu.emacs.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>gnu.emacs.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>exec "$HOME/Applications/Emacs.app/Contents/MacOS/Emacs" --fg-daemon --init-directory "$HOME/.config/emacs" --no-init-file --eval "(progn (setq epa-popup-info-window nil) (setq user-init-file (expand-file-name \"init.el\" user-emacs-directory)) (load user-init-file))"</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/emacs-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/emacs-daemon.log</string>
</dict>
</plist>
